# systemd/xhs-automation-worker-fast.service
[Unit]
Description=XHS Automation - Fast Worker (MCP Interactions)
After=network.target xhs-automation-main.service
Requires=xhs-automation-main.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/xhs-automation
ExecStart=/home/ubuntu/xhs-automation/venv/bin/huey_consumer.py src.queue.huey_config.huey -q fast -q normal -w 3
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/xhs-automation/logs/worker-fast.log
StandardError=append:/home/ubuntu/xhs-automation/logs/worker-fast.log
Environment="PYTHONPATH=/home/ubuntu/xhs-automation"

[Install]
WantedBy=multi-user.target
